# deploy/server/zaochang-backup.timer
[Unit]
Description=Daily Zaochang state backup

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m
Unit=zaochang-backup.service

[Install]
WantedBy=timers.target
